# 1h/1h.service
[Unit]
Description=1h cron

[Service]
User=abraxas
ExecStart=/home/abraxas/myq/service/1h.sh
